# serverhub-agent/packaging/debian/prerm
#!/bin/bash
set -e

case "$1" in
  remove|upgrade)
    systemctl stop serverhub-agent.service  || true
    systemctl disable serverhub-agent.service || true
    ;;
esac
